# srcs/7-grafana/docker-srcs/services.conf
[program:grafana]
command=/usr/sbin/grafana-server_wrapper.sh
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=5
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/proc/1/fd/2
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:sshd]
command=/usr/sbin/sshd -D -E /proc/1/fd/1
